# services/direwolf.timer
[Unit]
Description=Start direwolf on boot

[Timer]
AccuracySec=1
OnBootSec=80

[Install]
WantedBy=timers.target
