# services/multicasttnc0.timer
[Unit]
Description=Start multicast route for tnc0 on boot

[Timer]
AccuracySec=1
OnBootSec=124

[Install]
WantedBy=timers.target
