# services/socatPositRx.timer
[Unit]
Description=Start socat position Rx forwarding on boot

[Timer]
AccuracySec=1
OnBootSec=140

[Install]
WantedBy=timers.target
