# services/socatPositTx.timer
[Unit]
Description=Start socat position Tx forwarding on boot

[Timer]
AccuracySec=1
OnBootSec=120

[Install]
WantedBy=timers.target
