# services/socatChatTx.timer
[Unit]
Description=Start socat chat Tx forwarding on boot

[Timer]
AccuracySec=1
OnBootSec=115

[Install]
WantedBy=timers.target
