# services/tncattach.timer
[Unit]
Description=Start tncattach on boot

[Timer]
AccuracySec=1
OnBootSec=95

[Install]
WantedBy=timers.target
